# Example_Config_SYSTEMD/tdtrac-cron-EXAMPLE.service
[Unit]
Description=TDTrac CRON Queue Runner - EXAMPLE

[Service]
User=www-data
ExecStart=/path/to/installation/bin/cake cron_job
